# crontab
* * * * * /home/cary/.pyenv/versions/3.11.8/bin/gitpoll -p /home/cary/liberlife/gitpoll/pre.sh -P /home/cary/liberlife/gitpoll/post.sh /home/cary/liberlife