# systemd/liberlife.service
[Unit]
Description=LiberLife Service
After=network.target

[Service]
Type=simple
User=cary
WorkingDirectory=/home/cary/liberlife
Environment="LIBLIFE_ENV=prod"
EnvironmentFile=/home/cary/liberlife/.envs/prod.env
ExecStart=/bin/sh -c '. /home/cary/liberlife/shell.sh uvistart'
Restart=on-failure
[Install]
WantedBy=multi-user.target
